# setup/templates/keycloak/conf/startup.conf
[program:startup]
command=./root/salt-startup.sh
autostart=true
autorestart=false
user=root
startsecs=0
priority=50
