# setup/stack/states/base/mysql/mysql.conf
[program:mysqld]
command=/etc/supervisor/conf.d/mysql_starter.sh
autostart=true
autorestart=true
user=root

