# app/src/Setup/templates/mail/conf/salt-minion.conf
[program:salt-minion]
command=salt-minion
autostart=true
autorestart=false
user=root
