# app/src/Setup/templates/minion/salt-minion.conf
[program:salt-minion]
command=/bin/bash -c "salt-minion"
autostart=true
autorestart=false
user=root
