# app/src/Setup/stack/states/mailservices/mailservices/sv-apache.conf
[program:apache2]
command=/bin/bash -c "mkdir /var/run/apache2/;source /etc/apache2/envvars && exec /usr/sbin/apache2 -DFOREGROUND"
autostart=true
autorestart=true
user=root

