# app/src/Setup/stack/states/mailservices/mailservices/sv-postfix.conf
[program:postfix]
command=/root/postfix_fg_starter.sh
autostart=true
autorestart=false
user=root