# src/server/mail/conf/salt-minion.conf
[program:salt-minion]
command=/bin/bash -c "salt-minion -d"
autostart=true
autorestart=false
user=root
