# src/server/mail/conf/system.conf
[program:postfix]
command=/bin/bash -c "postfix start"
autostart=true
autorestart=false
user=root

[program:apache2]
command=/bin/bash -c "mkdir /var/run/apache2/;source /etc/apache2/envvars && exec /usr/sbin/apache2 -DFOREGROUND"
autostart=true
autorestart=true
user=root

[program:mysqld]
command=/usr/bin/pidproxy /var/run/mysqld/mysqld.pid /usr/bin/mysqld_safe --datadir=/var/lib/mysql --pid-file=/var/run/mysqld/mysqld.pid
autostart=true
autorestart=true
user=root
