# setup/stack/states/ilias/ilServer/ilServer.conf
[program:ilServer]
command=/bin/bash -c "/usr/bin/java -Dfile.encoding=UTF-8 -jar {{ ilserver_jar_path}} {{ ilserver_ini_path }} start"
autostart=true
autorestart=true
user=root
