# setup/stack/states/mailservices/mailservices/sv-dovecot.conf
[program:dovecot]
command=/bin/bash -c "dovecot"
autostart=true
autorestart=true
user=root

