# setup/stack/states/mailservices/mailservices/sv-mysql.conf
[program:mysqld]
command=/usr/bin/pidproxy /var/run/mysqld/mysqld.pid /usr/bin/mysqld_safe --datadir=/var/lib/mysql --pid-file=/var/run/mysqld/mysqld.pid
autostart=true
autorestart=true
user=root

