# setup/stack/states/mailservices/mailservices/sv-postfix.conf
[program:postfix]
command=/bin/bash -c "postfix start"
autostart=false
autorestart=false
user=root

