# src/lib/tpl/minion/salt-minion.conf
[program:salt-master]
command=/bin/bash -c "salt-minion -d"
autostart=true
autorestart=false
user=root
