# setup/stack/states/proxyservices/proxyservices/sv-nginx.conf
[program:nginx]
command=/etc/init.d/nginx start
autostart=true
autorestart=true
user=root

