# app/src/Setup/templates/keycloak/conf/mysql.conf
[program:mysqld]
command=/root/mysql_starter.sh
autostart=true
autorestart=true
user=root
priority=10
