# app/src/Setup/templates/keycloak/conf/salt-minion.conf
[program:salt-minion]
command=salt-minion
autostart=true
autorestart=true
user=root
priority=100
